# mount/tunnel.mount
[Unit]
Description=foo automount

[Automount]
Where=/mnt/foo

[Install]
WantedBy=multi-user.target
